# crontab
* * * * * root /app/clean-files.sh >> /var/log/cron.log 2>&1
